# crontab
0 7 * * * root /usr/local/bin/python3 /bots/veilleesr-bot.py --jorf  >> /var/log/cron.log  2>&1
0 8 * * 1 root /usr/local/bin/python3 /bots/veilleesr-bot.py --jorfrecap  >> /var/log/cron.log  2>&1
0 13 * * * root /usr/local/bin/python3 /bots/veilleesr-bot.py --retweet --tweet --datarand --jorf  >> /var/log/cron.log  2>&1
#
